# ops/finalrep-autodeploy.service
[Unit]
Description=FinalRep auto deploy check
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/finalrep
ExecStart=/opt/finalrep/deploy.auto.sh
